# roles/result/files/result.service
[Unit]
Description=Poll web application

[Service]
Environment
WorkingDirectory=/usr/local/bin/result/
ExecStart=npm start
Restart=always